# lab2/pkg/DEBIAN/postinst
#!/bin/sh
python -m compileall /bin/HTTPServerPythonSource/
update-rc.d -f task2.sh defaults
